# build/sni/etc/inittab
::sysinit:/scripts/init.sh
::askfirst:/bin/busybox sh

